# terminal-server.service
[Unit]
Description=Terminal WebSocket Server
After=network.target

[Service]
ExecStart=python3 /home/clungus/terminal/server.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
